# .UpdateBackups/rc.local
#!/bin/sh
mount -o bind /volume1/@optware /opt 
# Optware setup
[ -x /etc/rc.optware ] && /etc/rc.optware start

exit 0
